# chimera-proxy.socket
[Unit]
Description=Chimera socket

[Socket]
ListenStream=80
Service=chimera-proxy.service

[Install]
WantedBy=sockets.target
