# steam-buddy-proxy.socket
[Unit]
Description=SteamBuddy socket

[Socket]
ListenStream=80
Service=steam-buddy-proxy.service

[Install]
WantedBy=sockets.target
