# steam-patch.service
[Unit]
Description=Steam Patch Service

[Service]
ExecStart=/usr/bin/steam-patch
Restart=always

[Install]
WantedBy=default.target
