# steam-buddy.service
[Unit]
Description=SteamBuddy
After=network.target

[Service]
ProtectSystem=strict
ReadWritePaths=/var/lib/steam-buddy /tmp
WorkingDirectory=/usr/share/steam-buddy
ExecStart=/usr/share/steam-buddy/steam-buddy
Restart=on-abort
